# ansible/roles/cirismanager/templates/ciris-manager.service.j2
[Unit]
Description=CIRISManager - CIRIS Agent Orchestration Service
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User={{ manager_user }}
WorkingDirectory={{ manager_install_dir }}
Environment="PATH={{ manager_install_dir }}/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="CIRIS_MANAGER_CONFIG={{ manager_config_dir }}/config.yml"
{% if manager_dev_mode | default(false) %}
Environment="CIRIS_DEV_MODE=true"
{% endif %}
Environment="MANAGER_JWT_SECRET={{ manager_jwt_secret | default('dev-jwt-secret-key-for-testing-only') }}"
Environment="CIRIS_ENCRYPTION_SALT={{ manager_encryption_salt | default('ciris-manager-encryption-salt') }}"

ExecStart={{ manager_install_dir }}/venv/bin/ciris-manager --config {{ manager_config_dir }}/config.yml

Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
# ProtectHome disabled - manager needs ~/.local for Python packages
ProtectHome=false
ReadWritePaths={{ manager_data_dir }} {{ manager_log_dir }} {{ manager_config_dir }} /etc/nginx/sites-enabled
PrivateTmp=true
# Set HOME to writable location
Environment="HOME={{ manager_data_dir }}"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ manager_service_name }}

[Install]
WantedBy=multi-user.target
